# com.vx.market-alerts.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.vx.market-alerts</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/hang/Documents/Vx/scripts/start-market-alerts.sh</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/hang</string>
    <key>PATH</key>
    <string>/Users/hang/.nvm/versions/node/v22.17.1/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>WorkingDirectory</key>
  <string>/Users/hang/Documents/Vx</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/hang/Documents/Vx/market-alerts.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/hang/Documents/Vx/market-alerts.err.log</string>
</dict>
</plist>
